# root/etc/services.d/fail2ban/run
#!/usr/bin/with-contenv bash

  exec \
        fail2ban-client -x -vv -f start
